# deploy/macos/com.maki.llama-router.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.maki.llama-router</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__REPO_ROOT__/scripts/run-llama-router.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/.runtime/llama-router.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/.runtime/llama-router-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME__</string>
      <key>PATH</key>
      <string>__PATH_FOR_LAUNCHD__</string>
      <key>ROUTER_PORT</key>
      <string>8081</string>
    </dict>
  </dict>
</plist>
