# deploy/macos/com.maki.gateway.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.maki.gateway</string>
    <key>ProgramArguments</key>
    <array>
      <string>__NODE_BIN__</string>
      <string>__REPO_ROOT__/apps/gateway/server.mjs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/.runtime/gateway.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/.runtime/gateway-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>__HOME__</string>
      <key>PATH</key>
      <string>__PATH_FOR_LAUNCHD__</string>
      <key>LLAMA_SERVER_URL</key>
      <string>http://127.0.0.1:8081</string>
    </dict>
  </dict>
</plist>
